# Cron Jobs/reboot.cron
@reboot /bin/sh -l -c $SOHBASEDIR/PythonScripts/rebooted.py
